# deploy/systemd/gemini-coding-factory.service
[Unit]
Description=Gemini Coding Factory - Multi-Repository AI Development Service
Documentation=https://github.com/wjorgensen/Gemini-team
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=wes
Group=wes
WorkingDirectory=/home/wes/gemini-coding-factory
ExecStart=/usr/bin/node /home/wes/gemini-coding-factory/dist/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gemini-factory

# Environment variables
Environment=NODE_ENV=production
Environment=WORKSPACE_ROOT=/home/wes/coding-factory
Environment=PORT=3000
Environment=LOG_LEVEL=info

# Environment file for secrets
EnvironmentFile=-/etc/gemini-coding-factory/environment

# Security settings - Relaxed for development service that needs git/npm access
NoNewPrivileges=true
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/wes/coding-factory /home/wes/gemini-coding-factory /var/log/gemini-coding-factory
PrivateTmp=false
# Commented out restrictive settings that cause namespace issues
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 